# telegram_downloader.service
[Unit]
Description=Telegram Media Downloader
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/telegram_downloader
ExecStart=/usr/bin/python3 /path/to/telegram_downloader/telegram_downloader.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
